# php-worker/supervisord.conf
[supervisord]
nodaemon=true

[program:yike-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/yike/artisan queue:work redis --sleep=3 --tries=3
autostart=true
autorestart=true
user=root
numprocs=8
redirect_stderr=true
stdout_logfile=/var/www/yike/storage/logs/worker.log
